# scripts/te-broker/geard.service
[Unit]
Description=Geard deamon

[Service]
ExecStart=/bin/geard -d

