# scripts/te-broker/te_workers.service
[Unit]
Description=TE Workers

[Service]
ExecStart=/opt/stack/tripleo-ci/scripts/te-broker/start_workers.sh

